# contrib/init/bitcoinabc.conf
description "Bitcoin ABC Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BITCOINABC_BIN="/usr/bin/bitcoinabc"
env BITCOINABC_USER="bitcoin"
env BITCOINABC_GROUP="bitcoin"
env BITCOINABC_PIDDIR="/var/run/bitcoinabc"
# upstart can't handle variables constructed with other variables
env BITCOINABC_PIDFILE="/var/run/bitcoinabc/bitcoinabc.pid"
env BITCOINABC_CONFIGFILE="/etc/bitcoinabc/bitcoinabc.conf"
env BITCOINABC_DATADIR="/var/lib/bitcoinabc"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # bitcoinabc will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BITCOINABC_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run bitcoinabc."
        echo "The setting must appear in $BITCOINABC_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BITCOINABC_PIDDIR"
    chmod 0755 "$BITCOINABC_PIDDIR"
    chown $BITCOINABC_USER:$BITCOINABC_GROUP "$BITCOINABC_PIDDIR"
    chown $BITCOINABC_USER:$BITCOINABC_GROUP "$BITCOINABC_CONFIGFILE"
    chmod 0660 "$BITCOINABC_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BITCOINABC_PIDFILE" \
    --chuid $BITCOINABC_USER:$BITCOINABC_GROUP \
    --exec "$BITCOINABC_BIN" \
    -- \
    -pid="$BITCOINABC_PIDFILE" \
    -conf="$BITCOINABC_CONFIGFILE" \
    -datadir="$BITCOINABC_DATADIR" \
    -disablewallet \
    -daemon
